# app/autopilot.service
[Unit]
Description=Trading Bot Autopilot Daemon
After=network.target postgresql.service

[Service]
Type=simple
WorkingDirectory=/root/trading-bot/app
ExecStart=/usr/bin/python3 /root/trading-bot/app/autopilot_daemon.py
Restart=always
RestartSec=5
EnvironmentFile=/root/trading-bot/app/.env

[Install]
WantedBy=multi-user.target
