# app/live_test_2_report.timer
[Unit]
Description=Live Test 2 report timer — 2026-02-20 09:00 UTC (18:00 KST)

[Timer]
OnCalendar=2026-02-20 09:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
